# wombll.conf
description "Run gunicorn"
start on runlevel [2345]
stop on runlevel [06]
respawn
respawn limit 10 5
chdir /vagrant
exec /usr/local/bin/gunicorn wombll.wsgi
